# pop-up/init.service
[Unit]
Description=Ver usuarios

[Service]
ExecStart=/home/xd/Escritorio/pop-up/gui.sh

[Install]
WantedBy=multi-user.target